# packaging/debian/prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/emperor" ]; then
  if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
    invoke-rc.d emperor stop || true
  else
    /etc/init.d/emperor stop || true
  fi
fi